# misc/mcp.service
[Unit]
Description=DAC Master Control Program
After=mpd.service

[Service]
Type=forking
User=root
ExecStart=/usr/local/bin/mcp -d

[Install]
WantedBy=multi-user.target
