# misc/mcp.service
[Unit]
Description=DAC Master Control Program
After=mosquitto.service mpd.service network-online.target

[Service]
Type=simple
User=hje
StateDirectory=mcp
RuntimeDirectory=mcp
RuntimeDirectoryPreserve=yes
ExecStart=/usr/local/bin/mcp
#ExecStartPre=/usr/bin/ir-keytable -c -w /etc/rc_keymaps/universum

[Install]
WantedBy=multi-user.target
