# misc/mcp.service
[Unit]
Description=DAC Master Control Program
After=mpd.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/mcp

[Install]
WantedBy=multi-user.target
